# myapp.conf
[program:myapp]
priority=10
directory=/tmp/flask-example
command=python app.py
user=root
autostart=true
autorestart=true
